# plexmediaserver.conf
[program:plexmediaserver]
user = root
command = bash /opt/plexmediaserver/start.sh
